# packaging/omnios/boxvault-smf.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    SMF manifest for BoxVault Vagrant Box Repository Management System
-->
<service_bundle type='manifest' name='boxvault'>
  <service name='system/virtualization/boxvault' type='service' version='1'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    
    <!-- Dependencies -->
    <dependency name='filesystems' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    
    <!-- Start method -->
    <exec_method name='start' type='method' exec='/opt/boxvault/startup.sh' timeout_seconds='60'>
      <method_context working_directory='/opt/boxvault'>
        <method_credential user='boxvault' group='boxvault' 
            privileges='basic,net_privaddr,file_dac_read' />
        <method_environment>
          <envvar name='PATH' value='/opt/ooce/bin:/opt/ooce/node-22/bin:/usr/gnu/bin:/usr/bin:/usr/sbin:/sbin'/>
          <envvar name='NODE_ENV' value='production'/>
          <envvar name='CONFIG_PATH' value='/etc/boxvault'/>
          <envvar name='HOME' value='/var/lib/boxvault'/>
        </method_environment>
      </method_context>
    </exec_method>
    
    <!-- Stop method -->
    <exec_method name='stop' type='method' exec='/opt/boxvault/shutdown.sh' timeout_seconds='30'>
      <method_context working_directory='/opt/boxvault'>
        <method_credential user='boxvault' group='boxvault' />
        <method_environment>
          <envvar name='PATH' value='/opt/ooce/bin:/opt/ooce/node-22/bin:/usr/bin:/usr/sbin:/sbin'/>
        </method_environment>
      </method_context>
    </exec_method>
    
    <!-- Refresh method -->
    <exec_method name='refresh' type='method' exec=':kill -HUP' timeout_seconds='30' />
    
    <!-- Stability -->
    <stability value='Evolving' />
    
    <!-- Service description -->
    <template>
      <common_name>
        <loctext xml:lang='C'>BoxVault Vagrant Box Repository Management System</loctext>
      </common_name>
      <description>
        <loctext xml:lang='C'>
          Comprehensive Vagrant box repository management system that provides
          a web interface for managing, organizing, and distributing Vagrant
          boxes across teams and organizations with multi-organization support,
          role-based access control, and RESTful API for automation.
        </loctext>
      </description>
      <documentation>
        <manpage title='boxvault' section='8' />
        <doc_link name='BoxVault Project'
                  uri='https://github.com/Makr91/BoxVault' />
      </documentation>
    </template>
  </service>
</service_bundle>
