# packaging/DEBIAN/systemd/boxvault.service
[Unit]
Description=BoxVault Vagrant Box Repository Management System
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=boxvault
Group=boxvault
WorkingDirectory=/opt/boxvault
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=CONFIG_DIR=/etc/boxvault

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/boxvault /var/log/boxvault /etc/boxvault /local
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boxvault

# Resource limits
MemoryMax=12G
TasksMax=600

[Install]
WantedBy=multi-user.target
